# backup/crontab
0 * * * * cd /usr/local/my-finance/ && /bin/bash run.sh
